# packaging/pler.service
[Unit]
Description=pler - eBPF execve monitor for Wazuh
Documentation=https://github.com/pler/pler
After=network.target
ConditionPathExists=/sys/kernel/btf/vmlinux

[Service]
Type=simple
ExecStart=/usr/bin/pler --config /etc/pler/config.yaml
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pler

User=root
# kernel >= 5.8: CAP_BPF + CAP_PERFMON cukup
# kernel 4.18 (RHEL/Rocky 8): butuh CAP_SYS_ADMIN untuk perf mmap
AmbientCapabilities=CAP_BPF CAP_PERFMON CAP_SYS_RESOURCE CAP_SYS_ADMIN
CapabilityBoundingSet=CAP_BPF CAP_PERFMON CAP_SYS_RESOURCE CAP_SYS_ADMIN

PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/log/pler /run/pler /etc/pler
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
