# packaging/preremove.sh
#!/bin/sh
systemctl stop pler || true
systemctl disable pler || true
